# deploy/systemd/smokemon.service
[Unit]
Description=smokemon node agent (sample, detect, ship, prune)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
EnvironmentFile=-__ENV_FILE__
ExecStart=__PYTHON__ __AGENT__ collect
Restart=always
RestartSec=5
# The agent must never be the reason a box is slow. It yields to everything else and its own
# writes go to the back of the IO queue -- on an SD card that matters more than the CPU nice.
Nice=10
IOSchedulingClass=idle
# It reads /proc and /sys, writes one database, and talks to the hub. Nothing else.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__DATA_DIR__
PrivateTmp=true
NoNewPrivileges=true
# ping needs a raw ICMP socket. The kernel does not apply a binary's file capability
# (setcap cap_net_raw+ep on fping) through execve when the calling process has
# NoNewPrivileges=true -- fping runs, opens no socket, and reports total loss on every
# target with no error, which reads as a genuine outage everywhere at once. Ambient
# capabilities are the directive built for exactly this: granted to this unit's own
# process and inherited by fping across exec, unaffected by NoNewPrivileges. install.sh's
# setcap step still runs, for anyone invoking fping outside this unit (`smoke`, by hand).
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
