# deploy/systemd/smokemon-shipper.service
[Unit]
Description=smokemon shipper (push deltas to hub, one-shot drain)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__SMOKEMON_DIR__
EnvironmentFile=-__ENV_FILE__
ExecStart=__PYTHON__ -m smokemon.ship
