# deploy/systemd/smokemon-iperf.timer
[Unit]
Description=smokemon iperf3 every 15 min

[Timer]
OnBootSec=120
OnUnitActiveSec=900
Persistent=true

[Install]
WantedBy=timers.target
