# deploy/launchd/com.smokemon.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.smokemon.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string><string>smokemon.cli</string><string>daily</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/YOUR_USERNAME/smokemon</string>
    <key>EnvironmentVariables</key>
    <dict><key>PYTHONPATH</key><string>/Users/YOUR_USERNAME/smokemon</string></dict>
    <key>StartCalendarInterval</key>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>55</integer></dict>
    <key>StandardOutPath</key><string>/Users/YOUR_USERNAME/smokemon/logs/daily.out.log</string>
    <key>StandardErrorPath</key><string>/Users/YOUR_USERNAME/smokemon/logs/daily.err.log</string>
</dict>
</plist>
