# deploy/launchd/com.smokemon.shipper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Optional: ship this node's rows to a central hub. Set the real secret first. -->
<plist version="1.0">
<dict>
    <key>Label</key><string>com.smokemon.shipper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string><string>smokemon.ship</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/YOUR_USERNAME/smokemon</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key><string>/Users/YOUR_USERNAME/smokemon</string>
        <key>SMOKEMON_HUB_URL</key><string>http://HUB-HOST:8765/ingest</string>
        <key>SMOKEMON_HUB_SECRET</key><string>changeme</string>
    </dict>
    <key>StartInterval</key><integer>60</integer>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>/Users/YOUR_USERNAME/smokemon/logs/shipper.out.log</string>
    <key>StandardErrorPath</key><string>/Users/YOUR_USERNAME/smokemon/logs/shipper.err.log</string>
</dict>
</plist>
