# launchd/com.stefan.smokemon-shipper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- smokemon shipper på macOS: skeppar delta till hubben var 60:e sek (en-shots-läge).
     Fyll i SMOKEMON_HUB_SECRET nedan så den matchar hubben. Kopiera till
     ~/Library/LaunchAgents/ och bootstrap:a som ovan. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stefan.smokemon-shipper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/stefan/smokemon/shipper.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/stefan/smokemon</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SMOKEMON_HUB_URL</key>
        <string>http://100.87.219.2:8765/ingest</string>
        <key>SMOKEMON_HUB_SECRET</key>
        <string>CHANGE_ME</string>
    </dict>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/stefan/smokemon/logs/shipper.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/stefan/smokemon/logs/shipper.err.log</string>
</dict>
</plist>
